# scripts/lab-baseline-cron.timer
[Unit]
Description=Run Dokuru Lab Baseline upload processor every 30 seconds

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=1s
Unit=lab-baseline-cron.service

[Install]
WantedBy=timers.target
